# templates/db2fmcd.service.j2
[Unit]
Description=DB2V105

[Service]
ExecStart={{db2_install.file}}/bin/db2fmcd
Restart=always

[Install]
WantedBy=default.target
